# feeds/packages/net/ddns-scripts/files/etc/init.d/.svn/text-base/ddns.svn-base
#!/bin/sh /etc/rc.common

START=95

start () {

   sh /usr/lib/ddns/dynamic_dns_updater.sh start 
}

stop () {

  echo "ddns stop"
}

