# target/linux/ramips/base-files/etc/hotplug.d/.svn/entries
10

dir
560
svn://svn.openwrt.org.cn/dreambox/trunk/target/linux/ramips/base-files/etc/hotplug.d
svn://svn.openwrt.org.cn/dreambox



2012-10-09T01:53:44.847730Z
554
zjhzzyf














74e81aba-36f0-4891-8cfb-1815f9025a0c

button
dir

firmware
dir

