# feeds/dreambox_packages/new-dialer/files/etc/init.d/ghca
#!/bin/sh /etc/rc.common
START=99

start (){
    /bin/re-dial &
}

stop (){
    killall re-dial
}
