# feeds/luci/luci/luci/applications/luci-aria2/root/etc/uci-defaults/luci-aria2
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@aria2[-1]
	set ucitrack.aria2="aria2"
	set ucitrack.aria2.init='aria2'
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
