# target/linux/brcm63x8/base-files/etc/hotplug.d/.svn/entries
10

dir
560
svn://svn.openwrt.org.cn/dreambox/trunk/target/linux/brcm63x8/base-files/etc/hotplug.d
svn://svn.openwrt.org.cn/dreambox



2012-08-13T15:46:25.941620Z
537
lintel














74e81aba-36f0-4891-8cfb-1815f9025a0c

button
dir

